# passenger.conf
[program:passenger]
user=root
command=/opt/start_passenger
redirect_stderr=true
autorestart=true
